# deploy/gmj-frontend.container
[Container]
Image=localhost/gmj-frontend:latest
ContainerName=gmj-frontend
Network=gmj.network
PublishPort=3000:3000

[Install]
WantedBy=default.target
